# luci-theme-footstrap/root/etc/uci-defaults/30_luci-theme-footstrap
#!/bin/sh

# Remove any previously registered footstrap theme entries. Older versions
# registered six dark/light variants — mode (dark/light) and palette are now
# client-side toggles (Appearance popover), so only the layout needs a
# server-side theme entry.
for t in Footstrap FootstrapDark FootstrapLight FootstrapTop FootstrapTopDark \
         FootstrapTopLight FootstrapSidebar FootstrapOnTop; do
	uci -q delete "luci.themes.$t"
done

# Register the two layouts.
uci set luci.themes.FootstrapSidebar=/luci-static/footstrap
uci set luci.themes.FootstrapOnTop=/luci-static/footstrap-top

# Normalise the active theme: collapse any footstrap dark/light variant onto its
# base layout (so upgrading users keep their layout), and default a fresh
# install to the sidebar layout. Non-footstrap themes are left untouched.
cur=$(uci -q get luci.main.mediaurlbase)
case "$cur" in
	/luci-static/footstrap-top*) uci set luci.main.mediaurlbase=/luci-static/footstrap-top ;;
	/luci-static/footstrap*)     uci set luci.main.mediaurlbase=/luci-static/footstrap ;;
	"") [ "$PKG_UPGRADE" != 1 ] && uci set luci.main.mediaurlbase=/luci-static/footstrap ;;
esac

uci commit luci
exit 0
