# luci-theme-footstrap/root/etc/uci-defaults/30_luci-theme-footstrap
#!/bin/sh

# HAVE WE BEEN INSTALLED BEFORE?
#
# The rule is: a FRESH install may make its sidebar layout the active theme; an
# UPGRADE must never change whatever theme the user is currently using. That was
# expressed as `[ "$PKG_UPGRADE" = 1 ]`, and in production PKG_UPGRADE is never set
# — apk does not export it and the package's postinst does not either, so only
# dev-sync.sh (which sets it by hand) ever took the upgrade branch. On a real router
# the guard was dead, and an upgrade that happened to find an empty mediaurlbase
# would activate this theme behind the user's back.
#
# A marker file on the rootfs is the state that actually distinguishes the two:
# it exists if and only if this script has completed at least once.
FS_MARKER=/usr/share/luci-theme-footstrap/.installed
if [ "$PKG_UPGRADE" = 1 ] || [ -f "$FS_MARKER" ]; then
	PKG_UPGRADE=1
else
	PKG_UPGRADE=0
fi

# Remove any previously registered footstrap theme entries, including the two-layout
# pair this release replaces. Mode, palette AND now the layout are all client-side
# toggles (Appearance popover), so exactly one server-side theme entry is left — it is
# re-added below.
for t in Footstrap FootstrapDark FootstrapLight FootstrapTop FootstrapTopDark \
         FootstrapTopLight FootstrapSidebar FootstrapOnTop; do
	uci -q delete "luci.themes.$t"
done

# Register the ONE theme. Sidebar and top bar are no longer two themes: they are one
# template whose chrome is morphed by :root[data-layout], toggled client-side in the
# Appearance popover. So the Design page lists a single "Footstrap" entry and the
# layout is not a server choice any more.
uci set luci.themes.Footstrap=/luci-static/footstrap

# Normalise the active theme. Every pre-consolidation media path below is neither
# registered in luci.themes nor present on disk any more: an install left pointing at
# one is stranded on a theme the Design page does not list, so it cannot be switched
# away from through the UI, and its media dir 404s. They ALL now map onto the single
# surviving path. Non-footstrap themes are left untouched, and a fresh install
# defaults to footstrap.
#
# MIGRATION — a router that was on the top-nav theme must STAY on the top bar.
# The layout now lives in the browser (localStorage 'fs-layout'), and a shell script
# cannot write that. So carry the fact server-side instead: footstrap_layout is the
# router's DEFAULT layout, which head.ut stamps onto <html data-layout> and the
# client's own choice then overrides. Set it only when the old path says the user was
# on top-nav, and only if it is not already set, so a re-run never clobbers a
# deliberate change.
cur=$(uci -q get luci.main.mediaurlbase)
new=$cur
case "$cur" in
	# was on the top-nav theme (incl. its legacy dark/light variants) -> keep the bar
	/luci-static/footstrap-top*)
		[ -n "$(uci -q get luci.main.footstrap_layout)" ] || \
			uci set luci.main.footstrap_layout=top
		new=/luci-static/footstrap ;;
	# any other footstrap path (incl. the legacy -dark/-light) -> the sidebar default
	/luci-static/footstrap*)          new=/luci-static/footstrap ;;
	"") [ "$PKG_UPGRADE" = 1 ] || new=/luci-static/footstrap ;;
esac

# Last-resort guard against a blank UI: whatever stays active must actually be on
# disk — both its media dir and its ucode template. If a partial upgrade or a
# hand-edited mediaurlbase points at neither, fall back to the stock bootstrap
# theme; a plain UI beats no UI. Skipped when bootstrap itself is not installed.
theme=${new##*/}
if [ -n "$theme" ] && { [ ! -d "/www/luci-static/$theme" ] ||
                        [ ! -f "/usr/share/ucode/luci/template/themes/$theme/header.ut" ]; }; then
	[ -d /www/luci-static/bootstrap ] && new=/luci-static/bootstrap
fi

[ "$new" = "$cur" ] || uci set luci.main.mediaurlbase="$new"

uci commit luci

# the dispatch/menu caches embed theme paths — drop them so the switch applies on
# the next request rather than after a reboot
rm -f /tmp/luci-indexcache* 2>/dev/null
rm -rf /tmp/luci-modulecache 2>/dev/null

# Drop the marker LAST, so a run that dies half-way is still treated as a fresh
# install next time rather than as an upgrade.
mkdir -p "$(dirname "$FS_MARKER")" 2>/dev/null
: > "$FS_MARKER" 2>/dev/null

exit 0
