# luci-theme-footstrap/root/etc/uci-defaults/30_luci-theme-footstrap
#!/bin/sh

changed=0

set_opt() {
	local key=$1
	local val=$2

	if ! uci -q get "luci.$key" 2>/dev/null; then
		uci set "luci.$key=$val"
		changed=1
	fi
}

set_opt themes.Footstrap /luci-static/footstrap

if [ "$PKG_UPGRADE" != 1 ] && [ $changed = 1 ]; then
	set_opt main.mediaurlbase /luci-static/footstrap
fi

set_opt themes.FootstrapDark /luci-static/footstrap-dark
set_opt themes.FootstrapLight /luci-static/footstrap-light

# top-nav variant (selectable in System → System → Language and Style)
set_opt themes.FootstrapTop /luci-static/footstrap-top
set_opt themes.FootstrapTopDark /luci-static/footstrap-top-dark
set_opt themes.FootstrapTopLight /luci-static/footstrap-top-light

if [ $changed = 1 ]; then
	uci commit luci
fi

exit 0
